# fab/files/ubuntu/etc/profile.d/fnet.sh
export ftp_proxy=http://hydrangeaproxy:1982166878@rep.proxy.nic.fujitsu.com:8080
export http_proxy=http://hydrangeaproxy:1982166878@rep.proxy.nic.fujitsu.com:8080
export https_proxy=http://hydrangeaproxy:1982166878@rep.proxy.nic.fujitsu.com:8080
